# deploy/systemd/orkestra-agent.service
[Unit]
Description=orkestra Agent — Docker Host Controller
Documentation=https://github.com/heckertobias/orkestra
After=docker.service network-online.target
Wants=docker.service network-online.target

[Service]
Type=simple
# The agent needs Docker socket access. Add the 'docker' group to SupplementaryGroups
# if running as a non-root user whose primary group cannot access the socket.
User=root
ExecStart=/usr/local/bin/orkestra-agent serve
EnvironmentFile=-/etc/orkestra/agent/env
Restart=on-failure
RestartSec=10s

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/orkestra /etc/orkestra/agent /run/docker.sock /var/run/docker.sock
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictRealtime=true
LockPersonality=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=orkestra-agent

[Install]
WantedBy=multi-user.target
